# bounsic-back/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
directory=/app/src
command=gunicorn main:app --bind 0.0.0.0:8000
autostart=true
autorestart=true
stdout_logfile=/var/log/gunicorn.log
stderr_logfile=/var/log/gunicorn.err

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.err
